# files/opt/venus/service/simple-upnpd/run
#!/bin/sh
echo "*** starting simple-upnpd ***"
exec 2>&1
exec /bin/start-simple-upnpd
