# files/opt/venus/service/netmon/run
#!/bin/sh
echo "*** starting netmon ***"
exec 2>&1
exec /opt/victronenergy/netmon/netmon
