# files/opt/venus/service/dbus-ble-sensors/run
#!/bin/sh
echo "*** starting dbus-ble-sensors ***"
exec 2>&1
exec softlimit -d 100000000 -s 1000000 -a 100000000 /opt/victronenergy/dbus-ble-sensors/start-ble-sensors.sh
