# files/opt/venus/service/venus-platform/run
#!/bin/sh
echo "*** starting venus-platform ***"
exec 2>&1
exec softlimit -d 100000000 -s 1000000 -a 100000000 /opt/victronenergy/venus-platform/venus-platform
