# files/opt/venus/service/dbus-adc/run
#!/bin/sh
echo "*** starting dbus-adc ***"
exec 2>&1
exec softlimit -d 100000000 -s 1000000 -a 100000000 /opt/victronenergy/dbus-adc/start-adc.sh
