# files/opt/venus/service/nginx/run
#!/bin/sh
echo "*** starting nginx ***"
exec 2>&1
exec /usr/sbin/start-nginx.sh
