# files/opt/venus/service/llmnrd/run
#!/bin/sh
echo "*** starting llmnrd ***"
exec 2>&1
exec setuidgid nobody /usr/sbin/llmnrd -H venus -6
