# files/opt/venus/service/vrmlogger/run
#!/bin/sh
echo "*** starting vrmlogger ***"
exec 2>&1
exec softlimit -d 100000000 -s 1000000 -a 100000000 /opt/victronenergy/vrmlogger/vrmlogger.py
