# files/opt/venus/service/service-advertiser/run
#!/bin/sh
echo "*** starting service-advertiser ***"
exec 2>&1
exec /opt/victronenergy/service-advertiser/navico-announcement.sh 80 /app
