# files/opt/venus/service/dbus-tempsensor-relay/run
#!/bin/sh
echo "*** starting dbus-tempsensor-relay ***"
exec 2>&1
exec softlimit -d 100000000 -s 1000000 -a 100000000 /opt/victronenergy/dbus-tempsensor-relay/dbus_tempsensor_relay.py
