# files/opt/venus/service/ppp/run
#!/bin/sh
echo "*** starting ppp ***"
exec 2>&1
exec /usr/sbin/pppd call provider nodetach
